# deploy/wattkeeper-agent.service
[Unit]
Description=Wattkeeper node agent
After=network-online.target nut-server.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wattkeeper-agent --config-dir /etc/nut --listen :80
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target