# image/stage-wattkeeper/03-firstboot/files/etc/systemd/system/wattkeeper-firstboot.service
[Unit]
Description=Wattkeeper first-boot initialization
After=apply_os_config.service local-fs.target systemd-remount-fs.service
Before=wattkeeper-agent.service
ConditionPathExists=!/var/lib/wattkeeper/.firstboot-complete

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/wattkeeper-firstboot.sh

[Install]
WantedBy=multi-user.target